# Docker-files/crontab
* * * * * cd /app/connector && /usr/local/bin/python3 ./vectra_saas_pull.py >> /var/log/cron.log 2>&1
